# thermometer.service
# /etc/systemd/system/thermometer.service

[Unit]
Description=Thermometer
After=network.target

[Service]
WorkingDirectory=/path/to/repo
ExecStart=/usr/bin/python3 /path/to/repo/app.py
Restart=always

[Install]
WantedBy=multi-user.target
